# supervisor/3-php.conf
[program:php]
command = /usr/sbin/php-fpm7.2 -F
autostart = true
autorestart=true
priority=2